# apis/course-manager-api/supervisor/course-manager.conf
[program:courses_service_course_application_workflow]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=course_application_workflow --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/course_application_workflow.log

[program:courses_service_course_application_workflow_email]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=course_application_workflow_email --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/course_application_workflow_email.log

[program:courses_service_upload_course_images_to_cdn]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=upload_images_to_cdn --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/upload-images-to-cdn.log

[program:courses_service_delete_course_images_to_cdn]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=delete_images_to_cdn --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/delete-images-to-cdn.log

[program:process_uploaded_courses_list]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=process_uploaded_courses_list --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/process-uploaded-courses-list.log

[program:courses_service_push_course_to_search_engine]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=push_course_to_search_engine --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/push-course-to-search-engine.log


[program:courses_service_update_the_course_standard_fee_payable_usd]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=update_standard_fee_payable_usd --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/course-updates.log

[program:courses_service_generate_course_index_list]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=generate_course_index_list --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/generate-course-index_list.log

[program:courses_service_delete_course_index]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=delete_course_index --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/delete_course_index.log

[program:courses_service_push_event]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=push_event --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/push-events.log

[program:courses_service_remove_course_from_search_engine]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=remove_course_to_search_engine --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/remove_course_from_search_engine.log

[program:courses_service_push_course_to_internal_search_engine]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=push_course_to_internal_search_engine --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/push-course-to-internal-search-engine.log

[program:courses_service_push_course_to_key_phrase_extractor]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=push_course_to_key_phrase_extractor --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/push-course-to-key-phrase-extractor.log


[program:courses_service_push_selected_key_phrase_to_internal_search_engine]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/course-manager-api/artisan queue:work --queue=push_selected_key_phrase_to_internal_search_engine --sleep=3 --tries=2 --timeout=1800 --daemon
autostart=true
autorestart=true
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/course-manager-api/storage/logs/push-selected_key-phrase-to-internal-search-engine.log




